# contrib/arkd.service
[Unit]
Description=arkd - Ark Protocol Server (Rust)
Documentation=https://github.com/lobbyclawy/arkd-rs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arkd
Group=arkd
ExecStart=/usr/local/bin/arkd --config /etc/arkd/config.toml
Restart=on-failure
RestartSec=5s
Environment=RUST_LOG=info
EnvironmentFile=-/etc/arkd/env
WorkingDirectory=/var/lib/arkd

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/arkd
PrivateTmp=true

[Install]
WantedBy=multi-user.target
